# deploy/systemd/teknav-backup.service
[Unit]
Description=Backup Teknav PostgreSQL database
After=postgresql.service

[Service]
Type=oneshot
User=teknav
Group=teknav
EnvironmentFile=/etc/teknav/backend.env
ExecStart=/opt/teknav/current/scripts/systemd-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
